# scripts/systemd/forest-sentinel-dashboard.service
[Unit]
Description=Open Forest Sentinel dashboard
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=@USER@
WorkingDirectory=@APP_DIR@
EnvironmentFile=@APP_DIR@/.env
ExecStart=@APP_DIR@/scripts/serve_dashboard.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
