# scripts/systemd/forest-sentinel-prune.timer
[Unit]
Description=Apply the Open Forest Sentinel COG retention policy on a schedule

[Timer]
# Rendered by vm_setup.sh from PRUNE_SCHEDULE (default 02:30:00 UTC daily,
# before the 03:00 pipeline run, so the day's disk headroom is reclaimed
# first). Pruning is safe alongside a running pipeline: only files older than
# the active window (plus margin) are ever deleted.
OnCalendar=@PRUNE_SCHEDULE@
Persistent=true

[Install]
WantedBy=timers.target
