# scripts/systemd/forest-sentinel-dashboard.service
[Unit]
Description=Open Forest Sentinel dashboard
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=@USER@
WorkingDirectory=@APP_DIR@
EnvironmentFile=@APP_DIR@/.env
ExecStart=/home/@USER@/.local/bin/uv run uvicorn forest_sentinel.dashboard.app:app --host 0.0.0.0 --port @DASHBOARD_PORT@
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
