# scripts/systemd/forest-sentinel-pipeline.service
[Unit]
Description=Open Forest Sentinel optical-change pipeline (one run over a rolling window)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=@USER@
WorkingDirectory=@APP_DIR@
EnvironmentFile=@APP_DIR@/.env
ExecStart=@APP_DIR@/scripts/run_pipeline.sh
# Backstop for a wedged run (Type=oneshot otherwise waits forever); the storage
# layer also times out each stuck Earth Engine export batch. With checkpointed,
# resumable runs a timeout mid-backfill is safe: the next timer firing resumes
# where this run left off. Configure via PIPELINE_TIMEOUT in config/instance.env.
TimeoutStartSec=@PIPELINE_TIMEOUT@
